# etc/profile.d/colors.sh
export aycolorw="#222222, #222222, #b1b4b3"
export aycolorn="#222222, #b1b4b3, #222222, #005577, #b1b4b3"
export aycolora="#222222, #b1b4b3, #222222, #007763, #b1b4b3"
export aycoloru="#222222, #b1b4b3, #222222, #77003d, #b1b4b3"
